# deploy/ostrave-bot.service
# systemd-юнит бота «Остров» — Linux-аналог restart_bot.ps1.
# Установка: см. DEPLOY.md в корне проекта.
#
# Гарантирует ровно один экземпляр (защита от TelegramConflictError),
# авто-рестарт при падении и автозапуск после перезагрузки сервера.

[Unit]
Description=Ostrave Telegram bot (заявки на вечеринку «Остров»)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Отдельный непривилегированный пользователь (создание — см. DEPLOY.md).
User=ostrave
Group=ostrave
WorkingDirectory=/opt/ostrave_bot
ExecStart=/opt/ostrave_bot/.venv/bin/python bot.py
Restart=on-failure
RestartSec=5

# Бот сам пишет ротируемый bot.log в WorkingDirectory;
# stdout/stderr дополнительно попадают в journald (journalctl -u ostrave-bot).
StandardOutput=journal
StandardError=journal

# Минимальное укрепление: только чтение системы, запись — лишь в каталог бота
# (там живут leads.db, fsm.db и bot.log).
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/ostrave_bot

[Install]
WantedBy=multi-user.target
